# root/etc/services.d/uploader-gdrive/finish
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

# Copyright (c) 2019, PhysK
# All rights reserved.

echo "Killing Uploader"
kill "$(pgrep uploader.sh)"